# romm/rootfs/etc/services.d/romm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Romm
# Runs Romm application
# ==============================================================================

# Wait for MariaDB and Redis to be ready
bashio::log.info "Waiting for database services to be ready..."
sleep 10

# Wait for MariaDB to accept connections
while ! mysqladmin ping -h127.0.0.1 --silent 2>/dev/null; do
    bashio::log.info "Waiting for MariaDB..."
    sleep 2
done

bashio::log.info "MariaDB is ready"

# Wait for Redis
while ! redis-cli -h 127.0.0.1 ping 2>/dev/null | grep -q PONG; do
    bashio::log.info "Waiting for Redis..."
    sleep 2
done

bashio::log.info "Redis is ready"

# Get configuration
DB_PASSWORD=$(bashio::config 'db_password')
AUTH_SECRET=$(bashio::config 'auth_secret_key')
LOG_LEVEL=$(bashio::config 'log_level')

# Database configuration
export DB_HOST="127.0.0.1"
export DB_PORT="3306"
export DB_NAME="romm"
export DB_USER="romm"
export DB_PASSWD="${DB_PASSWORD}"

# Redis configuration
export REDIS_HOST="127.0.0.1"
export REDIS_PORT="6379"

# Romm configuration
export ROMM_AUTH_SECRET_KEY="${AUTH_SECRET}"
export ROMM_AUTH_ENABLED="true"
export ROMM_BASE_PATH=""
export LOGLEVEL="${LOG_LEVEL}"

# Storage paths
export ROMM_DB_PATH="/data/config/romm.db"
export LIBRARY_PATH="/share/romm/library"
export RESOURCES_PATH="/data/resources"
export ASSETS_PATH="/data/assets"

# Metadata provider configuration
if bashio::config.has_value 'igdb_client_id'; then
    export IGDB_CLIENT_ID=$(bashio::config 'igdb_client_id')
fi

if bashio::config.has_value 'igdb_client_secret'; then
    export IGDB_CLIENT_SECRET=$(bashio::config 'igdb_client_secret')
fi

if bashio::config.has_value 'screenscraper_user'; then
    export SCREENSCRAPER_USERNAME=$(bashio::config 'screenscraper_user')
fi

if bashio::config.has_value 'screenscraper_password'; then
    export SCREENSCRAPER_PASSWORD=$(bashio::config 'screenscraper_password')
fi

if bashio::config.has_value 'steamgriddb_api_key'; then
    export STEAMGRIDDB_API_KEY=$(bashio::config 'steamgriddb_api_key')
fi

if bashio::config.has_value 'retroachievements_api_key'; then
    export RETROACHIEVEMENTS_API_KEY=$(bashio::config 'retroachievements_api_key')
fi

if bashio::config.has_value 'mobygames_api_key'; then
    export MOBYGAMES_API_KEY=$(bashio::config 'mobygames_api_key')
fi

# Scheduled scanning
if bashio::config.true 'enable_scheduled_rescan'; then
    export ENABLE_SCHEDULED_RESCAN="true"
    export SCHEDULED_RESCAN_CRON=$(bashio::config 'scheduled_rescan_cron')
fi

export SCAN_TIMEOUT=$(bashio::config 'scan_timeout')

bashio::log.info "Starting Romm application..."

# Change to Romm directory and start the application
cd /opt/romm || exit 1

# Activate virtual environment and run with gunicorn
exec .venv/bin/gunicorn app.main:app \
    --bind 0.0.0.0:8080 \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --access-logfile - \
    --error-logfile - \
    --log-level "${LOG_LEVEL,,}"
