# romm/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting nginx on port ${ROMM_PORT}..."

# Run envsubst on nginx templates to substitute environment variables
for template in /etc/nginx/templates/*.conf.template; do
    if [ -f "$template" ]; then
        envsubst < "$template" > "/etc/nginx/conf.d/$(basename "$template" .template)"
        bashio::log.info "Processed nginx template: $(basename "$template")"
    fi
done

exec nginx -g 'daemon off; error_log /dev/stderr info;'
