# romm/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

# Get ingress port for nginx
ROMM_PORT=$(bashio::addon.ingress_port)
export ROMM_PORT="${ROMM_PORT:-8095}"
export ROMM_BASE_PATH="/romm"

bashio::log.info "Starting nginx on port ${ROMM_PORT}..."

# Run envsubst on nginx templates
for template in /etc/nginx/templates/*.conf.template; do
    if [ -f "$template" ]; then
        envsubst "${ROMM_PORT} ${ROMM_BASE_PATH}" < "$template" > "/etc/nginx/conf.d/$(basename "$template" .template)"
        bashio::log.info "Processed nginx template: $(basename "$template")"
    fi
done

exec nginx -g 'daemon off; error_log /dev/stderr info;'
