# profilarr-v1/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio

bashio::log.warning "nginx service stopped"
