# romm/rootfs/etc/services.d/mariadb/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Romm
# Handles MariaDB shutdown
# ==============================================================================

bashio::log.info "MariaDB service stopped."

# Cleanup PID file
rm -f /var/run/mysqld/mysqld.pid
