# cleanuparr/rootfs/etc/services.d/cleanuparr/finish
#!/usr/bin/with-contenv bashio
# =============================================================================
# Home Assistant Add-on: Cleanuparr
# Service Finish Script - Handles graceful shutdown
# =============================================================================

bashio::log.info "Cleanuparr service stopping..."

# Get the exit code from the service (first argument)
EXIT_CODE=${1:-0}

# If the service was killed by a signal, the exit code is 256
# and the signal number is in the second argument
if [ "${EXIT_CODE}" -eq 256 ]; then
    SIGNAL=${2:-TERM}
    EXIT_CODE=$((128 + SIGNAL))
    bashio::log.info "Cleanuparr killed by signal ${SIGNAL}"
else
    bashio::log.info "Cleanuparr exited with code ${EXIT_CODE}"
fi

# Write the exit code for the container
# This ensures the container exits with the correct status
echo "${EXIT_CODE}" > /run/s6-linux-init-container-results/exitcode

bashio::log.info "Cleanuparr service stopped with exit code ${EXIT_CODE}"
