# romm/rootfs/etc/services.d/redis/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Romm
# Handles Redis shutdown
# ==============================================================================

bashio::log.info "Redis service stopped."
