# romm/rootfs/etc/services.d/redis/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Romm
# Runs Redis
# ==============================================================================

bashio::log.info "Starting Redis..."

exec redis-server \
    --dir /data/redis \
    --dbfilename dump.rdb \
    --bind 127.0.0.1 \
    --port 6379 \
    --save 60 1 \
    --loglevel notice
