# .claude/skills/ha-docker-addons/examples/simple-wrapper/rootfs/etc/s6-overlay/s6-rc.d/init-config/up
#!/command/with-contenv bashio

exec 2>&1

bashio::log.info "Generating configuration..."

# Create data directory
mkdir -p /data/myapp

# Set environment variables for gomplate
export APP_HOST="0.0.0.0"
export APP_PORT=$(bashio::addon.port 8080)
export APP_LOG_LEVEL=$(bashio::config 'log_level')
export APP_RETENTION=$(bashio::config 'data_retention')

# Generate configuration from template
gomplate \
  -f /defaults/app-config.yaml.gotmpl \
  -o /data/myapp-config.yaml

if [ $? -ne 0 ]; then
  bashio::log.error "Configuration generation failed"
  exit 1
fi

bashio::log.info "Configuration generated successfully"
exit 0
