# .claude/skills/ha-docker-addons/examples/simple-wrapper/rootfs/etc/s6-overlay/s6-rc.d/myapp/run
#!/command/with-contenv bashio

exec 2>&1

bashio::log.info "Starting application..."

# Get configuration
LOG_LEVEL=$(bashio::config 'log_level')
PORT=$(bashio::addon.port 8080)
CONFIG_FILE="/data/myapp-config.yaml"

# Verify config exists
if [ ! -f "${CONFIG_FILE}" ]; then
  bashio::log.error "Configuration file not found: ${CONFIG_FILE}"
  exit 1
fi

bashio::log.info "Using log level: ${LOG_LEVEL}"
bashio::log.info "Using port: ${PORT}"
bashio::log.info "Configuration: ${CONFIG_FILE}"

# Run application (using busybox as placeholder)
# In real add-on, this would be: exec /usr/bin/myapp
exec /usr/bin/myapp httpd -f -p "${PORT}" -h /data/myapp
